# files-common/etc/crontabs/root
# VPN watchdog каждые 10 минут
*/10 * * * * /usr/bin/vpn_watchdog.sh
# DNS watchdog каждые 5 минут
*/5 * * * * /usr/bin/dns_watchdog.sh
# FRP watchdog каждые 5 минут (проверяет сервер и переключает если нужно)
*/5 * * * * /usr/bin/frpc_watchdog.sh
# Проверка обновлений каждый день в 3:00
0 3 * * * /usr/bin/update_all.sh all
# Перезагрузка раз в неделю в воскресенье 4:00
0 4 * * 0 reboot
